# dot_config/systemd/clear-downloads.service
[Unit]
Description=Remove old files in Downloads
