# dot_config/systemd/user/systemd-tmpfiles-clean.timer.d/override.conf
[Timer]
AccuracySec=6h
RandomizedDelaySec=1h
Persistent=true
