# dot_config/systemd/user/plasma-ksplash.service.d/override.conf
[Unit]
# Try to make sure the splashscreen failing doesn't result in systemd saying
# the system is "degraded"
SuccessExitStatus=TERM 0
CollectMode=inactive-or-failed
