# dot_config/autostart/vlc-speakers.desktop
[Desktop Entry]
Type=Service
Actions=open;
Exec=vlc ~/Music/keep_speakers_on.flac
Type=Application
Terminal=false
Icon=vlc
